# setup/systemd/flask.service
[Unit]
Description=Gunicorn instance to serve Flask app
After=network.target

[Service]
User=pi
Group=www-data
WorkingDirectory=/var/www/html/raspberryprojects/python/flask_project
Environment="PATH=/var/www/html/raspberryprojects/python/flask_project/venv/bin"
Environment="FLASK_ENV=production"
ExecStart=/var/www/html/raspberryprojects/python/flask_project/venv/bin/gunicorn --workers 1 --threads 8 --timeout 60 --bind 0.0.0.0:5000 manage:app
Restart=always
RestartSec=5s
Umask 007

[Install]
WantedBy=multi-user.target
